# continuous_integration/docker_hadoop/images/hadoop-testing-base/files/etc/master.supervisord.d/hdfs-namenode.conf
[program:hdfs-namenode]
command=hdfs namenode
startsecs=2
stopwaitsecs=10
user=hdfs
redirect_stderr=true
stdout_logfile=/var/log/hadoop-hdfs/hadoop-hdfs-namenode.log
autostart=true
autorestart=false
