# continuous_integration/docker_hadoop/images/hadoop-testing-base/files/etc/master.supervisord.d/hdfs-init.conf
[program:hdfs-init]
command=sudo -u hdfs /root/setup-hdfs.sh
startsecs=2
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/hdfs-init.log
autostart=true
autorestart=false
