# continuous_integration/docker_hadoop/files/etc/edge-kerb.supervisord.conf
[supervisord]
strip_ansi = true
nodaemon = true
logfile = /var/log/supervisord.log
pidfile = /var/run/supervisord.pid

[unix_http_server]
file = /tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
prompt = edge

[group:main]
programs=

[program:config-then-startup]
command=/root/use-kerb-config-then-start.sh
startsecs=0
stopwaitsecs=10
user=root
redirect_stderr=true
stdout_logfile=/var/log/config_then_startup.log
autostart=true
autorestart=false
